# infrastructure/provisions/letsencrypt.crontab
2 3 5 * * dehydrated --cron --domain `hostname`
